# debian/python3-moat-gpio.service
[Unit]
Description=Run moat-gpio
After=network.target

[Service]
ExecStart=/usr/bin/moat-gpio /etc/moat/gpio.cfg
